# setup/init/sysvinit/supervizio
#!/bin/sh
### BEGIN INIT INFO
# Provides:          supervizio
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: PID1-capable process supervisor
# Description:       supervizio process supervisor daemon
### END INIT INFO

NAME="supervizio"
DAEMON="/usr/local/bin/supervizio"
DAEMON_ARGS="--config /etc/supervizio/config.yaml"
PIDFILE="/var/run/${NAME}.pid"
LOGFILE="/var/log/supervizio/supervizio.log"

[ -x "$DAEMON" ] || exit 0

do_start() {
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "$NAME is already running"
        return 1
    fi
    echo "Starting $NAME..."
    mkdir -p /var/log/supervizio
    start-stop-daemon --start --background \
        --make-pidfile --pidfile "$PIDFILE" \
        --exec "$DAEMON" -- $DAEMON_ARGS
    echo "$NAME started"
}

do_stop() {
    if [ ! -f "$PIDFILE" ] || ! kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "$NAME is not running"
        return 1
    fi
    echo "Stopping $NAME..."
    start-stop-daemon --stop --pidfile "$PIDFILE" \
        --signal TERM --retry "TERM/30/KILL/5"
    rm -f "$PIDFILE"
    echo "$NAME stopped"
}

do_status() {
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "$NAME is running (PID $(cat "$PIDFILE"))"
        return 0
    fi
    echo "$NAME is not running"
    return 3
}

case "$1" in
    start)   do_start ;;
    stop)    do_stop ;;
    restart) do_stop; do_start ;;
    status)  do_status ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac
