# setup/init/systemd/supervizio.service
[Unit]
Description=Supervizio Process Supervisor
Documentation=https://github.com/kodflow/supervizio
After=network.target
Wants=network.target

[Service]
Type=notify
ExecStartPre=/bin/mkdir -p /etc/supervizio
ExecStartPre=/bin/chmod 755 /etc/supervizio
ExecStart=/usr/local/bin/supervizio --config /etc/supervizio/config.yaml
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Run as root for PID1 supervisor capabilities
User=root
Group=root

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/supervizio /var/run/supervizio /etc/supervizio
RuntimeDirectory=supervizio
LogsDirectory=supervizio

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=SUPERVIZIO_LOG_LEVEL=info

[Install]
WantedBy=multi-user.target
